# generators/security-defender-endpoint/samples/alert_profiles.json
[
  {
    "alert_id_prefix": "da638034938542563831",
    "title": "Suspicious credential access via LSASS",
    "category": "CredentialAccess",
    "severity": "High",
    "service_source": "Microsoft Defender for Endpoint",
    "detection_source": "EDR",
    "attack_techniques": "[\"OS Credential Dumping (T1003)\",\"LSASS Memory (T1003.001)\"]"
  },
  {
    "alert_id_prefix": "da638034938542563832",
    "title": "Suspicious PowerShell command line",
    "category": "Execution",
    "severity": "Medium",
    "service_source": "Microsoft Defender for Endpoint",
    "detection_source": "EDR",
    "attack_techniques": "[\"Command and Scripting Interpreter (T1059)\",\"PowerShell (T1059.001)\"]"
  },
  {
    "alert_id_prefix": "da638034938542563833",
    "title": "Persistence via registry Run key modification",
    "category": "Persistence",
    "severity": "Medium",
    "service_source": "Microsoft Defender for Endpoint",
    "detection_source": "EDR",
    "attack_techniques": "[\"Boot or Logon Autostart Execution (T1547)\",\"Registry Run Keys (T1547.001)\"]"
  },
  {
    "alert_id_prefix": "da638034938542563834",
    "title": "Suspicious network connection to known malicious IP",
    "category": "InitialAccess",
    "severity": "High",
    "service_source": "Microsoft Defender for Endpoint",
    "detection_source": "EDR",
    "attack_techniques": "[\"Phishing (T1566)\",\"Spearphishing Link (T1566.002)\"]"
  },
  {
    "alert_id_prefix": "da638034938542563835",
    "title": "Defense evasion via process injection",
    "category": "DefenseEvasion",
    "severity": "High",
    "service_source": "Microsoft Defender for Endpoint",
    "detection_source": "EDR",
    "attack_techniques": "[\"Process Injection (T1055)\",\"Dynamic-link Library Injection (T1055.001)\"]"
  },
  {
    "alert_id_prefix": "da638034938542563836",
    "title": "Lateral movement via PsExec",
    "category": "LateralMovement",
    "severity": "Medium",
    "service_source": "Microsoft Defender for Endpoint",
    "detection_source": "EDR",
    "attack_techniques": "[\"Remote Services (T1021)\",\"SMB/Windows Admin Shares (T1021.002)\"]"
  },
  {
    "alert_id_prefix": "da638034938542563837",
    "title": "Discovery activity: network share enumeration",
    "category": "Discovery",
    "severity": "Low",
    "service_source": "Microsoft Defender for Endpoint",
    "detection_source": "EDR",
    "attack_techniques": "[\"Network Share Discovery (T1135)\"]"
  },
  {
    "alert_id_prefix": "da638034938542563838",
    "title": "Ransomware behavior detected",
    "category": "Impact",
    "severity": "High",
    "service_source": "Microsoft Defender for Endpoint",
    "detection_source": "EDR",
    "attack_techniques": "[\"Data Encrypted for Impact (T1486)\"]"
  },
  {
    "alert_id_prefix": "da638034938542563839",
    "title": "Scheduled task created for persistence",
    "category": "Persistence",
    "severity": "Low",
    "service_source": "Microsoft Defender for Endpoint",
    "detection_source": "EDR",
    "attack_techniques": "[\"Scheduled Task/Job (T1053)\",\"Scheduled Task (T1053.005)\"]"
  },
  {
    "alert_id_prefix": "da638034938542563840",
    "title": "Antivirus detection of malware",
    "category": "Execution",
    "severity": "Informational",
    "service_source": "Microsoft Defender for Endpoint",
    "detection_source": "Antivirus",
    "attack_techniques": "[\"User Execution (T1204)\"]"
  },
  {
    "alert_id_prefix": "da638034938542563841",
    "title": "Suspicious use of certutil for file download",
    "category": "DefenseEvasion",
    "severity": "Medium",
    "service_source": "Microsoft Defender for Endpoint",
    "detection_source": "EDR",
    "attack_techniques": "[\"Ingress Tool Transfer (T1105)\",\"Deobfuscate/Decode Files or Information (T1140)\"]"
  },
  {
    "alert_id_prefix": "da638034938542563842",
    "title": "Brute force authentication attempt",
    "category": "CredentialAccess",
    "severity": "Medium",
    "service_source": "Microsoft Defender for Endpoint",
    "detection_source": "EDR",
    "attack_techniques": "[\"Brute Force (T1110)\",\"Password Spraying (T1110.003)\"]"
  }
]
